# upscaled-tl/com.upscaled.tl-sheets-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.upscaled.tl-sheets-sync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>/Users/connorodea/Library/Mobile Documents/com~apple~CloudDocs/UPSCALED2026/upscaled-tl/run_tl_sheets_sync.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>172800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/connorodea/Library/Logs/upscaled-tl-sheets-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/connorodea/Library/Logs/upscaled-tl-sheets-sync.err.log</string>
  </dict>
</plist>
